# systemd/hiddify-bot.service
[Unit]
Description=Hiddify VPN Bot v4.0 with PostgreSQL and Redis
Documentation=https://github.com/hiddify/hiddify-manager
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=hiddify
Group=hiddify
WorkingDirectory=/opt/hiddify-manager

# Environment
Environment="PATH=/opt/hiddify-manager/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/opt/hiddify-manager/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hiddify-bot

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hiddify-manager/data /opt/hiddify-manager/logs

# Execution
ExecStart=/opt/hiddify-manager/venv/bin/python /opt/hiddify-manager/scripts/monitor_bot.py

# Graceful shutdown
TimeoutStopSec=30
KillSignal=SIGTERM
KillMode=mixed

[Install]
WantedBy=multi-user.target
